# packaging/debian-systemd/msmtpq-ng-smtpd.postinst
#!/bin/sh

set -e

systemctl enable sendmail-bs.socket
