# packaging/debian-systemd/msmtpq-ng-smtpd.prerm
#!/bin/sh

systemctl disable sendmail-bs.socket
